# supervisor/c9.conf
[program:c9]
command=node server.js -a : -w /usr/local/develop --port 9999 --listen 0.0.0.0
directory=/usr/local/c9/core
autostart=true
autorestart=true
